# config/abuseio_delegation.conf
[program:abuseio_queue_delegation]
command=php artisan queue:work --daemon --tries=5 --sleep=3 --memory=256 --delay=30 --queue=abuseio_delegation
directory=/abuseio
user=ubuntu
