# config/abuseio_email_incoming.conf
[program:abuseio_queue_email_incoming]
command=php artisan queue:work --daemon --tries=1 --sleep=3 --memory=256 --delay=0 --queue=abuseio_email_incoming
directory=/abuseio
user=ubuntu
