# config/abuseio_email_outgoing.conf
[program:abuseio_queue_email_outgoing]
command=php artisan queue:work --daemon --tries=1 --sleep=3 --memory=256 --delay=0 --queue=abuseio_email_outgoing
directory=/abuseio
user=vagrant
