# config/crontab
# m h dom mon dow user  command
10 * * * *     root    /abuseio/bin/housekeeper
*/15 * * * *   root    /abuseio/bin/notifier
*/5 * * * *    root    /usr/bin/fetchmail &> /dev/null
